# wunzee.conf
[program:wunzee]
user=pi
command=python3 /home/pi/wunzee/wunzee.py
autostart=true
autorestart=true
