# wunzee.conf
[program:wunzee]
user=pi
command=python /home/pi/wunzee/wunzee.py
autostart=true
autorestart=true
